# infra/systemd/ocdash-worker@.service
[Unit]
Description=OpenCode Dashboard Worker (%i)
After=network.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__REPO_DIR__
EnvironmentFile=__REPO_DIR__/.env
Environment=OC_DASH_MODE=worker
Environment=OC_DASH_WORKER_ID=%i
ExecStart=/usr/bin/npm -w @ocdash/worker run dev
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
